# internal/installer/assets/slurm-k8s-kubelet.service
[Unit]
Description=Kubelet service

[Service]
ExecStart={{ .Scripts.Kubelet }}
Restart=on-failure
LimitNOFILE=65536